# services/unclutter/run
#!/bin/sh

exec unclutter --timeout 1 2>&1
