# xbindkeys/.config/systemd/user/xbindkeys.service
[Service]
ExecStart=xbindkeys -p -n
Restart=always

[Install]
WantedBy=default.target
