# xbindkeys/.config/systemd/user/xbindkeys.service
[Service]
ExecStart=xbindkeys -p -n

[Install]
WantedBy=default.target
