# crontab
*/15 * * * * . ~/.env && export PATH=~/.bin:$PATH && update-weather
