# crontab
0 */3 * * * newsbeuter -x reload
0   5 * * * . $HOME/.bashrc && perlbrew switch perl-5.10.1 && minicpan
