# crontab
*/15 * * * * /bin/bash -c 'export PATH=~/.bin/hush:~/.bin:/usr/bin:/usr/bin/vendor_perl:$PATH && update-weather'
