# crontab
*/15 * * * * /bin/bash -c'. ~/.env && export PATH=~/.bin:$PATH && update-weather'
