# crontab
*/15 * * * * /bin/bash -c 'export PATH=~/.bin:/bin:/usr/bin:/usr/bin/vendor_perl:$PATH && . ~/.env && update-weather'
