# mbsyncloop/.config/systemd/user/mbsyncloop.service
[Service]
ExecStart=%h/.bin/mbsyncloop
Restart=always

[Install]
WantedBy=default.target
