# ra-multiplex/.config/systemd/user/ra-multiplex.service
[Service]
ExecStart=%h/.cargo/bin/ra-multiplex-server
Restart=always

[Install]
WantedBy=default.target
