# i3status/.config/systemd/user/weather.timer
[Timer]
OnStartupSec=0
OnUnitActiveSec=900

[Install]
WantedBy=timers.target
