# i3status/.config/systemd/user/weather.service
[Service]
Type=oneshot
ExecStart=%h/.bin/%H/update-weather
