# xcape/.config/systemd/user/xcape.service
[Service]
ExecStart=xcape -d -e 'Control_L=Escape;Control_R=Escape'
Restart=always

[Install]
WantedBy=default.target
