# aranet/.config/systemd/user/aranet.service
[Service]
Type=oneshot
ExecStart=%h/.bin/sync-aranet E0:35:1F:BF:96:8D %h/Sync/aranet.db
