# aranet/.config/systemd/user/aranet.timer
[Timer]
OnStartupSec=0
OnUnitActiveSec=3600

[Install]
WantedBy=timers.target
