# services/kuake/run
#!/bin/sh

cd
exec urxvt -pe "kuake<XF86Tools>" -g 80x24 -name urxvt-kuake
