# services/kuake/run
#!/bin/sh

cd
exec urxvt -pe "kuake<F2>" -fn 7x14 -g 80x24 -name urxvt-kuake
