# services/kuake/run
#!/bin/sh

cd
urxvt -pe "kuake<F2>" -fn 6x10 -g 80x24 -name urxvt-kuake
