# services/compton/run
#!/bin/sh

exec compton -m 0.9 -i 0.8 -f -D 7 --backend glx --vsync drm
