# docker/.config/systemd/user/docker-prune.service
[Service]
Type=oneshot
ExecStart=docker system prune --filter until=168h --force
