# local/.config/systemd/user/watch-kbd.service
[Service]
ExecStart=%h/.bin/%H/watch-kbd

[Install]
WantedBy=default.target
