# services/xscreensaver/run
#!/bin/sh

xscreensaver -no-splash 2>&1
