# services/xscreensaver/run
#!/bin/sh

exec xscreensaver -no-splash 2>&1
