# services/xbindkeys/run
#!/bin/sh

cd $HOME
exec xbindkeys -p -n
