# laptop/.config/systemd/user/watch-kbd.service
[Service]
ExecStart=%h/.bin/watch-kbd
Restart=always

[Install]
WantedBy=default.target
