# server/.config/systemd/user/logout.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=sudo -K

[Install]
WantedBy=default.target
