# blueman/.config/systemd/user/blueman-applet.service.d/override.conf
[Service]
Restart=always

[Install]
WantedBy=default.target
