# services/minicpan_webserver/run
#!/bin/sh

cd $HOME
source $HOME/.bashrc
minicpan_webserver 2>&1
