# redshift/.config/systemd/user/redshift.service
[Service]
ExecStart=%h/.bin/redshift
Restart=always

[Install]
WantedBy=default.target
