# redshift/.config/systemd/user/redshift.service
[Service]
ExecStart=%h/.bin/redshift

[Install]
WantedBy=default.target
