# picom/.config/systemd/user/picom.service
[Service]
ExecStart=picom --config %h/.config/picom/picom.conf
Restart=always

[Install]
WantedBy=default.target
