# unclutter/.config/systemd/user/unclutter.service
[Service]
ExecStart=unclutter --timeout 1

[Install]
WantedBy=default.target
